# templates/supervisord.conf.j2
[supervisord]
user=root
nodaemon=true
directory={{ SUPERVISORLOGDIR }}
childlogdir={{ SUPERVISORLOGDIR }}
logfile=supervisord.log

[program:omero]
command={{ SUPERVISOROMEROSCRIPT }}
autostart=true
autorestart=true
startsecs=180
priority=100